# conf/systemd/dbs-worker.service
[Unit]
Description=DBS Worker

[Service]
ExecStart=/usr/bin/dbs celery worker -l INFO

[Install]
WantedBy=multi-user.target

